# systemd/fp-lapse.service
[Unit]
Description=fp-lapse — Sigma fp timelapse controller
After=multi-user.target
# The panel framebuffer may be fb0 (Bookworm) or fb1 (Trixie); the app
# locates it by driver name, so just require that some framebuffer exists.
ConditionPathExistsGlob=/dev/fb*
# Stop trying after 5 restarts within 60 s — a persistent crash loop
# shouldn't hammer the camera or burn battery.
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/home/pi/fp-lapse
# Local-only HTTP control surface for autonomous testing (see
# src/fp_lapse/control_server.py). Bound to 127.0.0.1 — never exposed
# off-box. Cheap to leave on; flip to 0 if you ever want to disable it.
Environment=FP_LAPSE_CONTROL=1
# The service runs as root but the venv is owned by the deploy user;
# don't write root-owned .pyc into it (that would later block the
# user's `uv sync`).
Environment=PYTHONDONTWRITEBYTECODE=1
ExecStart=/home/pi/fp-lapse/.venv/bin/python -m fp_lapse
# `always` (vs `on-failure`) also restarts if the process exits with
# status 0 — defensive: a clean exit in the field is itself a bug.
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# /dev/fb0 and GPIO need privileges; running as root is the path of least
# resistance for the field deployment. Tighten later with udev rules + groups
# (video, gpio, dialout) if desired.
User=root

[Install]
WantedBy=multi-user.target
